# debian/php-f8r-app.php-fpm.upstart
description "PHP FastCGI Process Manager"

start on (net-device-up and local-filesystems)
stop on runlevel [016]

limit nofile 5120 5120

expect fork
respawn
exec /content/bin/php-fpm -c /etc/php5/cgi/php.ini -y /etc/php5/cgi/php-fpm.conf
